# mcp-server/supervisor.conf
; Supervisor config for the taskbot MCP server (HTTP transport).
;
; The server listens only on 127.0.0.1:8765 — never directly exposed.
; Cloudflare Tunnel fronts it and provides the public HTTPS URL.
;
; Auth token is loaded from mcp-server/auth-token.txt by run.sh, so it
; doesn't appear in plaintext in this config.
;
; Install:
;   sudo cp supervisor.conf /etc/supervisor/conf.d/taskbot-mcp.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update

[program:taskbot-mcp]
directory=/home/telegrambot/projects/taskbot/mcp-server
command=/home/telegrambot/projects/taskbot/mcp-server/run.sh
user=telegrambot
autostart=true
autorestart=true
startsecs=3
startretries=10
stderr_logfile=/home/telegrambot/projects/taskbot/mcp-server/logs/supervisor.err.log
stdout_logfile=/home/telegrambot/projects/taskbot/mcp-server/logs/supervisor.out.log
environment=MCP_TRANSPORT="http",MCP_HOST="127.0.0.1",MCP_PORT="8765"
